# ping-monitor-daemon.service
[Unit]
Description=Ping Monitor daemon
After=graphical-session.target network-online.target
Wants=network-online.target
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=@@REPO_DIR@@
ExecStart=/usr/bin/env python3 @@REPO_DIR@@/ping-monitor-daemon.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
